# backend/db/crontab
00 00 * * * bash /home/backup.sh > /proc/1/fd/1 2> /proc/1/fd/2
# remember to end this file with an empty new line

